# Server/run-nodejs-server.ini
[program:run-nodejs-server]
directory = /home/argame/html/nodejs-server
command = node server.js
startsecs = 60
